# nginx-reload.path
[Path]
PathChanged=/etc/nginx-config-grug/

[Install]
WantedBy=multi-user.target
